# VAS/debian/postinst
#!/bin/sh
set -e

echo "[VAS] Ejecutando postinst..."

# Crear directorios persistentes
mkdir -p /var/lib/vas
mkdir -p /var/log/vas

# Crear DB si no existe
if [ ! -f /var/lib/vas/vas.db ]; then
    echo "[VAS] Creando base de datos inicial..."
    python3 /usr/lib/vas/database.py --init || true
fi

# Crear JSON inicial si no existe
if [ ! -f /var/lib/vas/computers.json ]; then
    echo '{"computers": []}' > /var/lib/vas/computers.json
fi

# Recargar systemd
systemctl daemon-reload || true

# Habilitar servicio
systemctl enable vas.service || true

# Arrancar en instalación fresca y reiniciar tras upgrade
if [ "$1" = "configure" ]; then
    if [ -n "${2:-}" ]; then
        systemctl restart vas.service || true
    else
        systemctl start vas.service || true
    fi
fi

echo "[VAS] Instalación completada."
