# debian/postinst
#!/bin/sh
set -e

echo "[VAS] Ejecutando postinst..."

# Crear directorios persistentes
mkdir -p /var/lib/vas
mkdir -p /var/log/vas
mkdir -p /etc/vas/vas.conf.d

# Crear JSON inicial si no existe
if [ ! -f /var/lib/vas/computers.json ]; then
    echo '{"computers": []}' > /var/lib/vas/computers.json
fi

# Crear archivo de versión si no existe
if [ ! -f /var/lib/vas/version ]; then
    echo "0" > /var/lib/vas/version
fi

# Recargar systemd
systemctl daemon-reload || true

# Habilitar servicio
systemctl enable vas.service || true

# Arrancar en instalación fresca y reiniciar tras upgrade
if [ "$1" = "configure" ]; then
    if [ -n "${2:-}" ]; then
        systemctl restart vas.service || true
    else
        systemctl start vas.service || true
    fi
fi

echo "[VAS] Instalación completada."

