# debian/postinst
#!/bin/sh
set -e

echo "[VAS] Ejecutando postinst..."

# Crear usuario de sistema si no existe
if ! getent passwd _vas > /dev/null 2>&1; then
    adduser --system --no-create-home --group _vas
fi

# Crear directorios persistentes con permisos correctos
mkdir -p /var/lib/vas
mkdir -p /var/log/vas
mkdir -p /etc/vas/vas.conf.d

chown -R _vas:_vas /var/lib/vas
chown -R _vas:_vas /var/log/vas

# Crear archivo de versión si no existe
if [ ! -f /var/lib/vas/version ]; then
    echo "0" > /var/lib/vas/version
    chown _vas:_vas /var/lib/vas/version
fi

# Recargar systemd
systemctl daemon-reload || true

# Habilitar servicio
systemctl enable vas.service || true

# Arrancar en instalación fresca y reiniciar tras upgrade
if [ "$1" = "configure" ]; then
    if [ -n "${2:-}" ]; then
        systemctl restart vas.service || true
    else
        systemctl start vas.service || true
    fi
fi

echo "[VAS] Instalación completada."
